# roles/bash_profile/files/etc/profile.d/local_bin.sh
umask 022
test -d ~/.local || mkdir -p ~/.local
test -d ~/.local/bin || mkdir -p ~/.local/bin
