# roles/bash_profile/files/etc/profile.d/umask.sh
if [ "$(umask)" = "0000" ]; then
    umask 022
fi
